# education_center.service
[Unit]
Description=Education Center Management System API
Documentation=https://github.com/yourusername/education_center_backend
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/education_center
Environment="PATH=/var/www/education_center/venv/bin"
Environment="PYTHONPATH=/var/www/education_center"
ExecStart=/var/www/education_center/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=yes
PrivateDevices=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/www/education_center/uploads
ReadWritePaths=/var/www/education_center/logs

[Install]
WantedBy=multi-user.target